# pkg/jispro/files-base/etc/profile.d/jispro-base.sh
# Managed by jispro-base package

export EDITOR='vim'
export SSH_AUTH_SOCK="$XDG_RUNTIME_DIR/ssh-agent.socket"
export PATH=$HOME/.local/bin:$PATH

alias ls=/usr/bin/exa

# export FZF_DEFAULT_COMMAND='rg --files'
# export FZF_DEFAULT_OPTS='-m --height 50% --border'